# files/dotfiles/.config/systemd/user/dropbox.service.d/override.conf
[Service]
Environment=DISPLAY=:0
ExecStartPre=/bin/sh -c '/usr/bin/kill -SIGINT $(pgrep -f "dropbox /newerversion")'
